# systemd/pacdb.path
[Path]
PathChanged=/var/lib/pacman/sync

[Install]
WantedBy=multi-user.target
