# systemd/pacdb.path
[Path]
PathModified=/var/lib/pacman/sync

[Install]
WantedBy=multi-user.target
