# systemd/pacdb.service
[Unit]
Description=pacdb Refresh Service

[Service]
Type=oneshot
ExecStart=/usr/bin/pacdb
